# deploy/groom.service
[Unit]
Description=Groom Agent
Documentation=https://github.com/etnz/groom
After=network-online.target
Wants=network-online.target
# Disable start rate limiting so systemd never gives up restarting the service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
# GoReleaser installs binaries to /usr/bin by default
ExecStart=/usr/bin/groom
# Fast restart to minimize downtime during updates
Restart=always
RestartSec=1
# Default configuration
Environment="GROOM_ADDR=:8080"

[Install]
WantedBy=multi-user.target